# debian/postinst
#!/bin/bash
# postinst script for flexibudget

set -e
. /usr/share/debconf/confmodule


#Replace Key Value in Config File
function replace {
    cfg="/usr/share/flexibudget/includes/config.php"
    sed -i "/${1}'/c\define('${1}', '${2}');" $cfg
}

#Replace Key Value in PhinX Config File
function replaceYML {
    cfg="/var/lib/flexibudget/phinx.yml"
    sed -i "/${1}'/c\        ${1}: ${2}" $cfg
}

listOfConfigs="FLEXIBEE_URL
FLEXIBEE_LOGIN
FLEXIBEE_PASSWORD
FLEXIBEE_COMPANY
DATABASE_TYPE
DATABASE_HOST
DATABASE_PORT
DATABASE_USER
DATABASE_PASSWORD
DATABASE_NAME"

IFS='                                                                                                                                                        
'                                                                                                                                                            
for keyword in $listOfConfigs                                                                                                                                
do                                                                                                                                                           
  db_get "flexibudget/$keyword"
  if [ "$RET" = "false" ]; then
    db_input critical "flexibudget/$keyword" || true
    db_go
    fi
 replace "$keyword" $RET
done

if [ -x /usr/sbin/a2enconf ] ; then
    a2enconf flexibudget
fi

if [ -f /etc/init.d/apache2 ] ; then
    if [ -x /usr/sbin/invoke-rc.d ]; then
        invoke-rc.d apache2 reload 3>/dev/null || true
    else
        /etc/init.d/apache2 reload 3>/dev/null || true
    fi
fi

if [ -f /etc/init.d/avahi ] ; then
    if [ -x /usr/sbin/invoke-rc.d ]; then
        invoke-rc.d avahi reload 3>/dev/null || true
    else
        /etc/init.d/avahi reload 3>/dev/null || true
    fi
fi


cd /usr/share/flexibudget
if [ -f composer.lock ] ; then
    rm -f composer.lock
fi
composer --no-dev -o install

db_get "flexibudget/DATABASE_TYPE"
replaceYML "adapter" $RET
db_get "flexibudget/DATABASE_HOST"
replaceYML "host" $RET
db_get "flexibudget/DATABASE_PORT"
replaceYML "port" $RET
db_get "flexibudget/DATABASE_USER"
replaceYML "user" $RET
db_get "flexibudget/DATABASE_PASSWORD"
replaceYML "pass" $RET
db_get "flexibudget/DATABASE_NAME"
replaceYML "name" $RET

phinx -n -c/var/lib/flexibudget/phinx.yml migrate

#DEBHELPER#

exit 0
